# lidarr/lidarr.service
[Unit]
Description=Lidarr
After=syslog.target network.target

[Service]
User=ReplaceMe
Group=ReplaceMe

Type=simple
ExecStart=/usr/bin/mono /opt/lidarr/Lidarr.exe -nobrowser
TimeoutStopSec=20
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
